# templates/freebsd/rc.d/pyload.erb
#!/bin/sh
#
# This file was generated by Chef for <%= node['fqdn'] %>.
# Do NOT modify this file by hand!

# $FreeBSD$
#
# PROVIDE: pyload
# REQUIRE: LOGIN NETWORK
# KEYWORD: shutdown

. /etc/rc.subr

PATH="$PATH:/usr/local/bin"

name="pyload"
desc="Downloadtool for One-Click-Hoster written in python."
rcvar=${name}_enable

load_rc_config $name
: ${pyload_enable:="YES"}
: ${pyload_daemon:=""}
: ${pyload_user:=""}
: ${pyload_group:=""}
: ${pyload_pid:=""}
: ${pyload_options:=""}

pidfile="${pyload_pidfile}"
command="${pyload_daemon}"
command_interpreter="python"
command_args="${pyload_options}"

run_rc_command "$1"
