# templates/freebsd/pyload.erb
#!/bin/sh

# $FreeBSD$
#
# PROVIDE: pyload
# REQUIRE: LOGIN NETWORK
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf to enable this
# service during boot; setting optional parameters it is possible to overwrite
# defaults from this script
#
#   pyload_enable="YES"
#   # Optional
#   pyload_daemon="<%= @install_dir %>/pyLoadCore.py"
#   pyload_user="<%= @user %>"
#   pyload_group="<%= @group %>"
#   pyload_pid="<%= @pid_file %>"
#   pyload_options="--daemon --configdir=<%= @config_dir %> --pidfile=<%= @pid_file %>"

. /etc/rc.subr

PATH="$PATH:/usr/local/bin"

name="pyload"
rcvar=${name}_enable

: ${pyload_enable:="YES"}
: ${pyload_daemon:="<%= @install_dir %>/pyLoadCore.py"}
: ${pyload_user:="<%= @user %>"}
: ${pyload_group:="<%= @group %>"}
: ${pyload_pid:="<%= @pid_file %>"}
: ${pyload_options:="--daemon --configdir=<%= @config_dir %> --pidfile=<%= @pid_file %>"}

command_interpreter="python"
command="${pyload_daemon}"
command_args="${pyload_options}"
command_full="${command_interpreter} ${command} ${command_args}"

start_cmd="su -m ${pyload_user} -c '${command_full}'"
stop_cmd="su -m ${pyload_user} -c '${command_full} --quit'"

load_rc_config $name
run_rc_command "$1"
