# templates/default/rhel/init.d/pyload.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          pyload
# Required-Start:    $local_fs $remote_fs $network
# Required-Stop:     $local_fs $remote_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Downloadtool for One-Click-Hoster written in python.
# Description:       Downloadtool for One-Click-Hoster written in python.
### END INIT INFO

NAME="pyload"
DAEMON="<%= @install_dir %>/pyLoadCore.py"
USER="<%= @user %>"
GROUP="<%= @group %>"
PID_FILE="<%= @pid_file %>"
LOCK_FILE="/var/lock/$NAME"
OPTIONS="--daemon --configdir=<%= @config_dir %> --pidfile=<%= @pid_file %>"

START_TIMEOUT=2
STOP_TIMEOUT=2

[ -x $DAEMON ] || exit 1

[ -e /etc/sysconfig/$NAME ] && . /etc/sysconfig/$NAME

. /etc/init.d/functions

start () {
    daemon --pidfile "$PID_FILE" --user "$USER" "$DAEMON $OPTIONS"
    RETVAL=$?
    [ $RETVAL -eq 0 ] && touch $LOCK_FILE
    return $RETVAL
}

stop () {
    killproc -p $PID_FILE $DAEMON
    RETVAL=$?
    [ $RETVAL -eq 0 ] && rm -f $PID_FILE && rm -f $LOCK_FILE
    return $RETVAL
}

restart () {
    stop || RETVAL=$?
    [ -n "$STOP_TIMEOUT" ] && sleep $STOP_TIMEOUT
    start || RETVAL=$?
    [ -n "$START_TIMEOUT" ] && sleep $START_TIMEOUT
    status || RETVAL=$?
    RETVAL=$?
    return $RETVAL
}

status () {
    # as /etc/init.d/functions#status is not working during testing phase ...
    #status -p $PID_FILE -l $LOCK_FILE $DAEMON
    #RETVAL=$?
    #return $RETVAL
    # ... we need different approach of getting process status

    [ ! -f "$PID_FILE" ] && return 1
    pid=`cat $PID_FILE`
    [ -z "$pid" ] && return 1
    [ ! -d /proc/$pid ] &&  return 1
    cmd=`cat /proc/$pid/cmdline | tr '\000' '\n' | awk 'NR==2'`
    [ "$cmd" != "$DAEMON" ] &&  return 1
    return 0
}

case "$1" in
    start)
        echo -n $"Starting $NAME: "
        start
        echo
        ;;
    stop)
        echo -n $"Stopping $NAME: "
        stop
        echo
        ;;
    restart)
        echo -n $"Restarting $NAME: "
        restart
        echo
        ;;
    status)
        echo -n $"Checking status of $NAME: "
        if status ; then
            echo -n "running"
        else
            echo -n "apparently not running"
        fi
        echo
        ;;
    *)
        echo "Usage: /etc/init.d/$NAME {start|stop|restart|status}"
        exit 2
        ;;
esac

exit $?
